# .deploy/uwsgi_nginx_supervisord.conf
[program:uwsgi-app]
command = BLUEWARE_CONFIG_FILE=blueware.ini blueware-admin run-program /usr/local/bin/uwsgi --ini /goge/.deploy/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx